# root/etc/s6-overlay/s6-rc.d/syncthing/run
#!/command/with-contenv sh

exec s6-setuidgid appuser /opt/syncthing/syncthing -no-browser -no-restart -gui-address=0.0.0.0:8384 -home=/config
